# systemd/bot_v13.service
[Unit]
Description=Bot v13
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/app/bot_v13
EnvironmentFile=/app/bot_v13/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /app/bot_v13/main.py


# Restart on crash — wait 10s between restarts
Restart=on-failure
RestartSec=10s

# Output to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot_v13

# Resource limits
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
